# drupal.cron
# cron for Drupal.
# bounce messages back to root
MAILTO=root

0 * * * * http cd /usr/share/drupal && /usr/bin/php4 -q cron.php

# vim:syn=crontab
